# server-config/systemd/cards.socket
[Unit]
Description=gunicorn socket for cards

[Socket]
ListenStream=/home/sandbox/cards/cards.sock
User=www-data

[Install]
WantedBy=sockets.target
